# Mule-Salt/supervisord.conf
[supervisord]
nodaemon=true

[program:saltstack]
command=/start.sh

[program:mule]
command=/opt/mule/bin/mule